# issbot@.service
[Unit]
Description=Telegram ISS bot

[Service]
User=%i
ExecStart=/usr/bin/python /home/rocky/issbot/issbot.py

[Install]
WantedBy=multi-user.target
